# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:frontend]
directory=/app/frontend
command=node server.js
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
directory=/app/backend
command=java -jar app.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SPRING_DATASOURCE_URL="jdbc:postgresql://${DB_HOST:-postgres}:${DB_PORT:-5432}/${DB_NAME:-career_nav}",SPRING_DATASOURCE_USERNAME="${POSTGRES_USER:-career_user}",SPRING_DATASOURCE_PASSWORD="${POSTGRES_PASSWORD:-career_pass}",AI_SERVICE_BASE_URL="http://127.0.0.1:8000"

[program:ai-service]
directory=/app/ai-service
command=/app/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=OLLAMA_BASE_URL="http://${OLLAMA_HOST:-ollama}:11434"
